# docker-monitor.service
[Unit]
Description=Docker Container Monitor
Documentation=https://github.com/rmnlu/container-monitor/
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=rda
Group=docker
WorkingDirectory=/data/monitoring/container-monitor
ExecStart=/data/monitoring/.venv/bin/python3 /data/monitoring/container-monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=docker-monitor

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/monitoring/container-monitor/logs
ReadOnlyPaths=/data/monitoring/container-monitor/docker_monitor

# Resource limits
CPUQuota=50%
MemoryLimit=512M

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
