# docker-monitor.service
[Unit]
Description=Docker Container Monitor - Enhanced
Documentation=https://github.com/yourorg/docker-monitor
After=network.target docker.service postgresql.service
Requires=docker.service
Wants=postgresql.service

[Service]
Type=simple
User=monitor
Group=docker
WorkingDirectory=/opt/docker_monitor
ExecStart=/usr/bin/python3 /opt/docker_monitor/docker_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=docker-monitor

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/docker_monitor/logs
ReadOnlyPaths=/opt/docker_monitor

# Resource limits
CPUQuota=50%
MemoryLimit=512M

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
